# AutomationZ-OS-Setup/systemd/automationz-updater.service
[Unit]
Description=AutomationZ Updater (git pull)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/automationz-os/updater/automationz-update.sh
